# flasher@.service
[Unit]
Description=Run the flash util on a flash drive

[Service]
ExecStart=/usr/bin/flasher %I
