# software/90-flasher.rules
ACTION=="add|remove", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", ENV{DEVTYPE}=="disk", RUN+="/usr/bin/systemctl restart flasher-trigger@%k"
