# software/flasher-trigger@.service
[Unit]
Description=Run the flash util on a flash drive

[Service]
Type=oneshot
ExecStart=/usr/bin/flasher %I
