# software/flasher-trigger@.service
[Unit]
Description=Run the flash util on a flash drive

[Service]
Type=oneshot
ExecStart=/usr/lib/flasher/main.sh %I
